# crontab
#run python script every minutes
1 * * * * python3 /app/main.py